# airootfs/etc/systemd/system/nganjo-cpu-performance.service
[Unit]
Description=Set CPU governor to performance
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for cpu in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo performance > $cpu 2>/dev/null || true; done'
ExecStart=/bin/sh -c 'echo 0 > /sys/devices/system/cpu/intel_pstate/no_turbo 2>/dev/null || true'
ExecStart=/bin/sh -c 'echo 1 > /sys/devices/system/cpu/cpufreq/boost 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
